# bash/profile.d/15_boxen
source /opt/boxen/env.sh

# brew-pip
export PYTHONPATH=$(brew --prefix)/lib/python2.7/site-packages
export PATH=$PATH:$(brew --prefix)/share/python
