# bash/profile.d/22_readline
# shell support for GNU Readline

export INPUTRC="${HOME}/.inputrc"
