# bash/profile.d/30_powerline
# Powerline

# export POWERLINE_COMMAND="${HOME}/Library/Python/2.7/bin/powerline"
export POWERLINE_COMMAND="/opt/boxen/pyenv/shims/powerline"

powerline_path=$(python -c 'import pkgutil; print pkgutil.get_loader("powerline").filename' 2>/dev/null)
# clobber existing prompt if we have Powerline
# if ! [[ -z "$powerline_path" ]]; then
#     source ${powerline_path}/bindings/bash/powerline.sh
# fi
