# bash/profile.d/40_vagrant
# Vagrant

# use Fusion if we have it installed
$(vagrant plugin list | grep -q fusion)
if [ $? ]; then
    VAGRANT_DEFAULT_PROVIDER=vmware_fusion
else
    VAGRANT_DEFAULT_PROVIDER=virtualbox
fi
export VAGRANT_DEFAULT_PROVIDER
