# bash/profile.d/44_anaconda
if [ -d ~/anaconda/bin ]; then
  export PATH="~/anaconda/bin:$PATH"
fi

# vim: set et fenc= ff=unix sts=2 sw=2 ts=2 : 
