# bash/profile.d/44_anaconda
if [ -d "${HOME}/anaconda/bin" ]; then
  export PATH="${HOME}/anaconda/bin:$PATH"
fi

# vim: set et fenc= ff=unix sts=2 sw=2 ts=2 : 
