# bash/profile.d/50_rk
# customizations for the rk user
#
if [[ "$USER" = 'rk' ]]; then
  # update global known hosts file
  curl -q -s --connect-timeout 1 -f 'http://puppet.runkeeper.com/statics/rk-known-hosts' -o '/tmp/rk-known-hosts' 2>/dev/null
  if [ $? -eq 0 ]; then
    mv /tmp/rk-known-hosts "${HOME}/.ssh/rk-known-hosts" >/dev/null 2>&1
  fi

  # patchHostgroup completion
  function _getopt_complete () {
    COMPREPLY=($( COMP_CWORD=$COMP_CWORD perl `which ${COMP_WORDS[0]}` ${COMP_WORDS[@]:0} ));
  }
  complete -F _getopt_complete patchHostgroup
fi
