# bash/profile.d/60_noti
# Noti
# https://github.com/variadico/noti
export NOTI_DEFAULT='pushover banner'
