# bash/profile.d/90_direnv
# direnv hooks
# shellcheck shell=bash

DIRENV_ACTIVATE="$(brew --prefix direnv)"

if [ -r "$DIRENV_ACTIVATE" ]; then
  eval "$(direnv hook bash)"
fi
